# t/fixtures/crontab
MAILTO="sample"

* * * * * ls -al

# blank lines