# t/fixtures/crontab

* * * * * ls -al

### start: remove me
# take me out
### end: remove me
# remains

# blank lines